# scripts/iot-hub-cv.service
[Unit]
Description=IoT Hub CV Pipeline (Hailo-8 NPU)
Documentation=https://github.com/your-org/iot-hub
# Start after Docker is up so mosquitto + mediamtx containers are available.
After=network-online.target docker.service local-fs.target
Wants=network-online.target
RequiresMountsFor=/mnt/ssd

[Service]
Type=simple
User=vlad
Group=vlad
WorkingDirectory=/mnt/ssd/iot-hub/repo

# Pipeline defaults — all can be overridden in /etc/iot-hub-cv.env
Environment=RTSP_URL=rtsp://localhost:8554/camera
Environment=HEF_PATH=/opt/iot-hub/models/current_yolo.hef
Environment=POSE_HEF_PATH=/opt/iot-hub/models/current_pose.hef
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
# ROOM is only the startup fallback — the pipeline polls BACKEND_URL
# /api/cv/pipeline-config and adopts the room of the camera's floor-plan
# placement, so moving the camera in the editor needs no change here.
Environment=ROOM=living_room
Environment=BACKEND_URL=http://localhost:8000
Environment=TARGET_FPS=15
Environment=METRICS_PORT=8002
Environment=ENABLE_FUSION=true

# Optional site-local overrides (create this file to change any of the above)
EnvironmentFile=-/etc/iot-hub-cv.env

# Uses a dedicated venv created with system Python 3.13 + --system-site-packages
# so hailo_platform (cpython-313) and cv2 (python3-opencv) are accessible.
# One-time setup:
#   python3 -m venv --system-site-packages /mnt/ssd/iot-hub/repo/.venv-cv
#   .venv-cv/bin/pip install -e /mnt/ssd/iot-hub/repo --no-deps
#   .venv-cv/bin/pip install aiomqtt prometheus-client pyyaml
ExecStart=/mnt/ssd/iot-hub/repo/.venv-cv/bin/python -m hub.edge.cv.pipeline

Restart=on-failure
RestartSec=10s
# Allow up to 5 restarts per 2 minutes before giving up
StartLimitIntervalSec=120
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=iot-hub-cv

# SIGTERM for graceful shutdown; pipeline closes Hailo device in finally block
KillSignal=SIGTERM
TimeoutStopSec=15

# TODO: ModelStore currently sends `docker kill --signal=SIGHUP cv` to trigger
# model reload. With this service, update ModelStore to use:
#   systemctl kill --signal=SIGHUP iot-hub-cv
# or write a PID file and use: kill -HUP $(cat /run/iot-hub-cv.pid)

[Install]
WantedBy=multi-user.target
