# scripts/model-puller.service
[Unit]
Description=IoT Hub model puller — pull new HEF artifacts from GHCR and promote
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
# Source the same .env that docker-compose uses so GHCR_TOKEN / DEPLOY_TOKEN
# are available without duplicating secrets.
EnvironmentFile=/opt/iot-hub/repo/.env
Environment=MODELS_DIR=/opt/iot-hub/models
Environment=BACKEND_URL=http://localhost:8000
ExecStart=/opt/iot-hub/repo/scripts/model-puller.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=model-puller
